# packer/files/reddit.service
[Unit]
Description=Reddit service

[Service]
Type=simple
ExecStart=/bin/bash /home/ubuntu/start-reddit.sh

[Install]
WantedBy=multi-user.target
